# etc/network/if-up.d/add-static-arp
#! /bin/sh

arp -i eth0 193.191.177.254 ca:fe:c0:ff:ee:00
